# yocto/raspberry_pi/bsp/sources/meta-ess/recipes-apps/flask-ess-application/files/flask-ess-application.service
[Unit]
Description=Flash ESS application

[Service]
Type=simple
Restart=always
RestartSec=10s
StandardOutput=journal+console
StandardError=journal+console
ExecStart=/usr/bin/flash-ess-application.py

[Install]
WantedBy=multi-user.target